# crontab/consumer.cron
# SYNTX Consumer - Batch Processing

# SYNTEX_SYSTEM Wrapper: 1x täglich (3 Uhr morgens, 20 Jobs)
0 3 * * * cd /opt/syntx-workflow-api-get-prompts && /usr/bin/python3 -c "import sys; sys.path.insert(0, 'queue_system'); from queue_system.core.consumer import QueueConsumer; c = QueueConsumer('syntex_system', 'cron_syntex'); stats = c.process_batch(20); print(f'Processed: {stats[\"processed\"]}, Failed: {stats[\"failed\"]}')" >> /opt/syntx-config/logs/consumer_syntex_cron.log 2>&1

# Sigma Wrapper: 4x täglich (4, 10, 16, 22 Uhr, 20 Jobs)
0 4,10,16,22 * * * cd /opt/syntx-workflow-api-get-prompts && /usr/bin/python3 -c "import sys; sys.path.insert(0, 'queue_system'); from queue_system.core.consumer import QueueConsumer; c = QueueConsumer('sigma', 'cron_sigma'); stats = c.process_batch(20); print(f'Processed: {stats[\"processed\"]}, Failed: {stats[\"failed\"]}')" >> /opt/syntx-config/logs/consumer_sigma_cron.log 2>&1
