# crontab/monitoring.cron
# SYNTX Monitoring - Stündliche Status-Logs
0 * * * * cd /opt/syntx-workflow-api-get-prompts && ./scripts/queue_status.sh >> /opt/syntx-config/logs/queue_status_hourly.log 2>&1

# Cleanup: Täglich um 2 Uhr (vor Consumer-Run)
0 2 * * * cd /opt/syntx-workflow-api-get-prompts && ./scripts/queue_cleanup.sh >> /opt/syntx-config/logs/cleanup.log 2>&1
