# crontab/producer.cron
# SYNTX Evolutionary Producer - Alle 2 Stunden
# Self-optimizing generation with field-based learning
0 */2 * * * cd /opt/syntx-workflow-api-get-prompts && /usr/bin/python3 evolution/evolutionary_producer.py >> /opt/syntx-config/logs/producer_cron.log 2>&1
